# infrastructure/crontab
0 0,12 * * * /usr/local/bin/php /app/bin/drm update_players >> /var/log/cron.log 2>&1
