# infrastructure/crontab
*/15 * * * * /usr/local/bin/php /app/bin/drm update_players >> /var/log/cron.log 2>&1
*/15 * * * * /usr/local/bin/php /app/bin/drm create_backup >> /var/log/cron.log 2>&1
